# ligo/cvmfs-oasis.opensciencegrid.org.automount

[Unit]
Documentation=Automount point for the OASIS CVMFS repository.
Before=local-fs.target

[Automount]
Where=/cvmfs/oasis.opensciencegrid.org

[Install]
WantedBy=local-fs.target
RequiredBy=local-fs.target

