# ligo/cvmfs-ligo.osgstorage.org.automount

[Unit]
Documentation=Automount point for the LIGO frame files CVMFS repository.
Before=local-fs.target

[Automount]
Where=/cvmfs/ligo.osgstorage.org

[Install]
WantedBy=local-fs.target
RequiredBy=local-fs.target

